# cadslicerprinter.service
[Unit]
Description=CadSlicerPrinter MCP Server
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Pc-MCP/src
Environment="WORKSPACE_DIR=/home/%i/Pc-MCP/workspace"
Environment="OPENSCAD_BIN=/usr/bin/openscad"
Environment="SLICER_BIN=/usr/bin/prusa-slicer"
Environment="OCTOPRINT_URL=http://localhost:5000"
Environment="OCTOPRINT_API_KEY="
# Default to localhost-only for security. Change --host to 0.0.0.0 for network access.
ExecStart=/usr/bin/python3 /home/%i/Pc-MCP/src/server.py --web --port 8080 --host 127.0.0.1
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
